# deploy/queue-worker.service
[Unit]
Description=HomeDuty Horizon
After=network.target mysql.service redis-server.service
Requires=mysql.service redis-server.service

[Service]
Type=simple
User=APP_USER
Group=www-data
WorkingDirectory=APP_DIR
ExecStart=/usr/bin/php artisan horizon
ExecStop=/usr/bin/php artisan horizon:terminate
Restart=always
RestartSec=5

StandardOutput=append:APP_DIR/storage/logs/horizon.log
StandardError=append:APP_DIR/storage/logs/horizon.log

[Install]
WantedBy=multi-user.target
